# systemd/xpr-oracle.service
[Unit]
Description=XPR Network oracle price pusher
Documentation=https://github.com/paulgnz/xpr-oracle
# nodeos is the push target. keosd is required for signing but is typically
# user-level (started ad-hoc by the same user that runs claimrewards), not a
# system service — so we don't require it as a systemd unit.
After=network-online.target nodeos.service
Wants=network-online.target
# Don't fork-bomb on a permanently bad config. NOTE: these belong in [Unit],
# not [Service] — systemd silently ignores them in [Service].
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# IMPORTANT: User/Group default to the same user that runs your claimrewards
# cron, so the daemon talks to the SAME keosd that already has your keys.
# install.sh sets these from --user=… (defaults to whoever runs install.sh).
# If you want a dedicated xpr-oracle system user with isolated keosd, see
# docs/BP-ONBOARDING.md §Advanced.
User=PLACEHOLDER_USER
Group=PLACEHOLDER_GROUP
WorkingDirectory=PLACEHOLDER_WORKING_DIR
Environment=XPR_ORACLE_CONFIG=PLACEHOLDER_CONFIG_PATH
# PATH for the cleos subprocess. The daemon shells out to `cleos` and inherits
# only what's in this unit's environment — there's no shell to read /etc/profile
# or ~/.bashrc. Common cleos install locations are below; if your BP's cleos
# lives somewhere else (e.g. /opt/eosio-2.0/bin), append it here.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/eosio/bin
# cleos looks for keosd's socket at $HOME/eosio-wallet/keosd.sock by default.
# install.sh sets HOME to the chosen user's actual home so cleos finds the
# right keosd without symlinking. If your keosd is on a different socket,
# set "walletUrl" in config.json instead.
Environment=HOME=PLACEHOLDER_HOME
# install.sh detects node via 'command -v node' and substitutes the absolute
# path here. nvm/fnm users may need to point this at their version-manager's
# node binary.
ExecStart=PLACEHOLDER_NODE_BIN dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
# SyslogIdentifier makes `journalctl -t xpr-oracle` work as a fallback on
# systems where `journalctl -u xpr-oracle.service` returns no entries (a
# known systemd 245 / Ubuntu 20.04 cgroup-attribution quirk with hardened
# units). On healthier systemd versions both `-u` and `-t` work.
SyslogIdentifier=xpr-oracle

# Hardening — note ProtectHome=read-only (not =true): the daemon runs as the
# same user that runs claimrewards, and it needs to talk to that user's keosd
# socket which lives in $HOME/eosio-wallet/. ProtectHome=true would render the
# whole /home tree empty for this unit; read-only keeps it visible. The keosd
# socket dir is added to ReadWritePaths below so cleos can write to the socket.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @debug
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
KeyringMode=private
DevicePolicy=closed

# Don't loop indefinitely on a permanently bad config
StartLimitIntervalSec=300
StartLimitBurst=5

# Allow writes to the keosd wallet directory so cleos can talk to keosd's
# Unix socket (the connection requires write access on the socket file).
# install.sh substitutes this with $HOME/eosio-wallet for the chosen user.
ReadWritePaths=PLACEHOLDER_KEOSD_DIR

[Install]
WantedBy=multi-user.target
